# docker/site/conf/celery.conf
[program:celery]
command=/home/ubuntu/site/bin/celery -A dmoj_celery worker
environment=DB_DATABASE=__DB_DATABASE__,DB_HOST=__DB_HOST__,DB_USER=__DB_USER__,DB_PASSWORD=__DB_PASSWORD__,DEBUG=__DEBUG__,TWO_STEP_REGISTRATION=__TWO_STEP_REGISTRATION__
directory=/home/ubuntu/site
# You should create a dedicated user for celery to run under.
user=ubuntu
group=ubuntu
stdout_logfile=/tmp/celery.stdout.log
stderr_logfile=/tmp/celery.stderr.log
